# wallchanger.cron
# m h  dom mon dow   command
*/5 * * * * gregnix /home/gregnix/.local/bin/wallchanger
