# wallchanger.cron
# m h  dom mon dow   command
0 * * * * gregnix /home/gregnix/.local/bin/wallchanger
